# syncthing-monitoring/syncthing-monitoring.timer
[Unit]
Description=Syncthing monitoring

[Timer]
OnBootSec=5m

[Install]
WantedBy=timers.target

