# auto_shutdown/shutdown_if_no_sessions.timer
[Unit]
Description=Shutdown if there are no tmux sessions

[Timer]
OnCalendar=*-*-* *:00:00

[Install]
WantedBy=timers.target
